# config/semaphore.service
[Unit]
Description=Semaphore UI
Documentation=https://semaphoreui.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=semaphore
Group=semaphore
WorkingDirectory=__DBA_PORTAL_ROOT__/semaphore
Environment=HOME=__DBA_PORTAL_ROOT__/home/semaphore
Environment=TMPDIR=__DBA_PORTAL_ROOT__/tmp
ExecStart=/usr/local/bin/semaphore server --config __DBA_PORTAL_ROOT__/semaphore/config.json
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=__DBA_PORTAL_ROOT__ /tmp

[Install]
WantedBy=multi-user.target
